# deploy/dodeka-backend-production.service
[Unit]
Description=Dodeka Backend Server - Production (Python)
After=network.target dodeka-auth-production.service
Wants=dodeka-auth-production.service

[Service]
Type=simple
User=backend
WorkingDirectory=/home/backend/dodeka/backend
ExecStart=/home/backend/.local/bin/uv run --frozen --no-dev production
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
