# deploy/dodeka-backend-demo.service
[Unit]
Description=Dodeka Backend Server - Demo (Python)
After=network.target dodeka-auth-demo.service
Wants=dodeka-auth-demo.service

[Service]
Type=simple
User=backend
WorkingDirectory=/home/backend/dodeka/backend
ExecStart=/home/backend/.local/bin/uv run --frozen --no-dev backend --env-file /home/backend/dodeka/backend/envs/demo/.env
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
