# deploy/dodeka-auth-demo.service
[Unit]
Description=Dodeka Auth Server - Demo (Go/tiauth-faroe)
After=network.target

[Service]
Type=simple
User=backend
WorkingDirectory=/home/backend/dodeka/backend/auth
ExecStart=/home/backend/dodeka/backend/auth/auth --env-file /home/backend/dodeka/backend/auth/envs/demo/.env
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
